# home/.config/autostart/albert.desktop
[Desktop Entry]
Type=Application
Name=Albert
Comment=Start albert at log in
Exec=albert

