# home/.config/autostart/compton.desktop
[Desktop Entry]
Type=Application
Name=Compton compositor
Comment=Start compton at log in
Exec=compton -cCG -d 5
OnlyShowIn=LXDE;
